# flexget.conf
# /etc/init/flexget.conf

description "FlexGet"
author "Alex Torrance"

start on opt-ready
stop on runlevel [06]

respawn
respawn limit 5 10

console log

setuid downloader

exec /opt/bin/flexget daemon start

pre-stop exec /opt/bin/flexget daemon stop
